# default-settings/usr/share/aurex-transmute/default-settings/base-files/etc/xdg/xdg-aurex/autostart/aurex-wallpaper-sync.desktop
[Desktop Entry]
Name=Wallpaper sync GNOME3 -> XFCE
Comment=Wallpaper sync GNOME3 -> XFCE
Exec=/usr/share/aurex-default-settigs/aurex-wallpaper-sync
Type=Application
OnlyShowIn=XFCE;
StartupNotify=false
