# default-settings/usr/share/aurex-transmute/default-settings/base-files/etc/xdg/xdg-aurex/autostart/aurex-kb-shortcuts.desktop
[Desktop Entry]
Name=aurex keyboard shortcuts updater
Comment=Update user's keyboard shortcurts
Exec=/usr/share/aurex-default-settings/aurex-kb-shortcuts
Type=Application
OnlyShowIn=XFCE;
StartupNotify=false
